# infra/aws-cli/market-copilot-frontend.service
[Unit]
Description=Market Copilot Frontend
After=network.target market-copilot-api.service

[Service]
User=marketcopilot
WorkingDirectory=/srv/market-copilot/frontend
EnvironmentFile=/etc/market-copilot/frontend.env
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
